# supervisord.conf
[supervisord]
nodaemon=true

[program:http-tunnel]
command=/usr/bin/autossh -qN http-tunnel

[program:nginx]
command=/docker/scripts/nginx-start
